# linkshortify-bot.service
[Unit]
Description=LinkShortify Telegram Bot
After=network.target mongod.service
Wants=mongod.service

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/linkshortify-master
ExecStart=/path/to/linkshortify-master/.venv/bin/python /path/to/linkshortify-master/bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=linkshortify-bot

[Install]
WantedBy=multi-user.target
